# roles/suttacentral/templates/sc-production-supervisor.conf
[program:sc-production]
command=/home/sc-production/suttacentral/daemonize
user=sc-production
autorestart=true
