# roles/suttacentral/templates/sc-staging-supervisor.conf
[program:sc-staging]
command=/home/sc-staging/suttacentral/daemonize
user=sc-staging
autorestart=true
